# units/systemd-detach-loopdev@.service.in
[Unit]
Description=Detach loop device %i
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=@rootlibexecdir@/systemd-detach-loopdev.sh %i
